# live_long_short.service
[Unit]
Description=prototype_long_short_v3 - Bybit Quant Bot (LIVE / live_trading.py)
After=network.target

[Service]
Type=simple
Environment=SSL_CERT_FILE=/usr/local/app/algo_3kings_long_short_v3/venv/lib/python3.12/site-packages/certifi/cacert.pem
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/usr/local/app/algo_3kings_long_short_v3
ExecStart=/usr/local/app/algo_3kings_long_short_v3/venv/bin/python3.12 -u live_trading.py
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
